# package/aws-iot/files/aws-iot.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org 

START=99
start() {
	if [ "$(cat /module_name)" = "RBR50" ]; then
		killall aws-iot
		/usr/sbin/aws-iot -t 86400 &	
	fi
}

stop() 
{
	killall aws-iot
}
boot()
{
	if [ "$(cat /module_name)" = "RBR50" ]; then
		/usr/sbin/aws-iot -t 86400 &	
	fi
}

